# ops/systemd/boomer-executor.service
[Unit]
Description=Boomer Executor — order management and Kite WebSocket tick feed
After=network.target

[Service]
Type=simple
User=boomer
Group=boomer
WorkingDirectory=/opt/boomer
EnvironmentFile=/etc/boomer/secrets.env
# Executor process manages the Kite WebSocket tick feed in-process (Q5-1 resolution).
# No separate boomer-websocket.service — the executor owns its own connection.
ExecStart=/opt/boomer/.venv/bin/python -m src.executor.main
Restart=always
RestartSec=10
StandardOutput=append:/var/log/boomer/executor.log
StandardError=append:/var/log/boomer/executor.log

[Install]
WantedBy=multi-user.target
